# systemd/i3status_touchy.service
[Unit]
# Not finished. idk

[Service]
ExecStart=/usr/local/bin/i3status_touchy

[Install]
WantedBy=default.target
